# deploy/cabal-smp.target
[Unit]
Description=Cabal SMP — Minecraft Java server
Requires=minecraft-cabal.service
After=network.target

[Install]
WantedBy=multi-user.target
