# deploy/minecraft-cabal.service
[Unit]
Description=Minecraft 26.2 server (minecraft-cabal)
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/minecraft-cabal/server
ExecStart=/root/minecraft-cabal/scripts/start.sh
ExecStop=/bin/kill -SIGINT $MAINPID
ExecStopPost=-/root/minecraft-cabal/scripts/capture-distance-manager-crash-context.sh
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
# Cap RSS for the whole service cgroup: matches scripts/start.sh heap (10G) plus off-heap / metaspace / native.
# If you change MIN_MEM/MAX_MEM in scripts/start.sh, raise MemoryMax accordingly (heap + ample headroom).
MemoryMax=20G

StandardInput=null
StandardOutput=journal
StandardError=journal
SyslogIdentifier=minecraft-cabal

[Install]
WantedBy=multi-user.target
